# other/battery-checker.service
[Unit]
Description=battery-checker system

[Service]
Type=simple
User=root
ExecStart=/bin/battery-checker

