# power/system/battery-checker.service
[Unit]
Description=suspends if low battery

[Service]
Type=oneshot
User=root
ExecStart=/bin/battery-checker

