# power/user/battery-checker.service
[Unit]
Description=suspends if low battery

[Service]
Type=oneshot
ExecStart=%h/.bin/battery-checker

