# xmonad/.config/systemd/user/failure-handler@.service
[Unit]
Description=Failure handler for %i

[Service]
Type=oneshot
ExecStart=notify-send 'Systemd unit failed' %i

