# clipmenu/.config/systemd/user/clipmenud.service.d/override.conf
[Unit]
PartOf=graphical-session.target

[Service]
Environment=CM_SELECTIONS=clipboard

[Install]
WantedBy=
WantedBy=graphical-session.target
