# themes/.config/systemd/user/feh.service
[Unit]
Description=Applies the last set background image

[Service]
Type=oneshot
ExecStart=%h/.fehbg

[Install]
WantedBy=graphical-session.target
